# monitor.service
[Unit]
Description=Monitor

[Service]
ExecStart=/root/monitor.sh

[Install]
WantedBy=multi-user.target
