# vmware-fix.service
[Unit]
Description=vmWare Fix

[Service]
ExecStart=/root/vmware-fix.sh

[Install]
WantedBy=multi-user.target
